# deploy/systemd/apexquant-v10.service
[Unit]
Description=ApexQuant V10 Paper Trading Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/apexquant-v10
EnvironmentFile=/var/www/apexquant-v10/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
